# installer/roundcube/roundcube-login-logger.service
[Unit]
Description=Roundcube login logger
Before=mysqld.target php-fpm.target nginx.target

[Service]
ExecStart=/usr/local/bin/roundcube-login-logger.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
